# debian/gamemode.postinst
#!/bin/sh
set -e

case "$1" in
	configure)
		# Disable the job by default for users updating from a buggy version
		if dpkg --compare-versions "$2" lt-nl 1.5.1-3; then
			if deb-systemd-helper --user --quiet was-enabled 'gamemoded.service'; then
				printf "Disable user unit that shouldn't have been auto enabled\n"
				deb-systemd-helper --user disable 'gamemoded.service' >/dev/null || true
			fi
		fi
		;;

	*)
	;;
esac

#DEBHELPER#
